# deploy/edge.service.example
[Unit]
Description=OLDFish.Claw Edge
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/oldfish-claw
EnvironmentFile=/opt/oldfish-claw/.env
ExecStart=/usr/bin/env npm run start:edge
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

